# dnsmasq-full/rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Preparing to start..."
bashio::config.require 'data_path'

DNSMASQ_DATA_PATH=$(bashio::config 'data_path')

#Create config if not exist
if ! bashio::fs.file_exists "$DNSMASQ_DATA_PATH/dnsmasq.conf"; then
  mkdir -p "$DNSMASQ_DATA_PATH" || bashio::exit.nok "Could not create $DNSMASQ_DATA_PATH"

  bashio::log.info "Create default config file"
  cat <<EOF > "$DNSMASQ_DATA_PATH/dnsmasq.conf"
no-resolv
no-hosts
keep-in-foreground
log-facility=-
no-poll
user=root
server=8.8.8.8
server=8.8.4.4
EOF
fi

# Run dnsmasq
bashio::log.info "Starting dnsmasq..."

# Set max open file limit to speed up startup
ulimit -n 1024

exec dnsmasq -C "$DNSMASQ_DATA_PATH/dnsmasq.conf" -z < /dev/null
