# dnsmasq-full/rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Preparing to start..."
bashio::config.require 'data_path'

DNSMASQ_DATA_PATH=$(bashio::config 'data_path')
CONFIG_FILE="$DNSMASQ_DATA_PATH/dnsmasq.conf"

#Create config if not exist
if ! bashio::fs.file_exists "$CONFIG_FILE"; then
  mkdir -p "$DNSMASQ_DATA_PATH" || bashio::exit.nok "Could not create $DNSMASQ_DATA_PATH"

  bashio::log.info "Create default config file"
  cat <<EOF > "$CONFIG_FILE"
no-resolv
no-hosts
keep-in-foreground
log-facility=-
no-poll
user=root
server=8.8.8.8
server=8.8.4.4
EOF
fi

# Run dnsmasq
bashio::log.info "Starting dnsmasq..."

# Set max open file limit to speed up startup
ulimit -n 1024

config_test_result=$(dnsmasq --test 2>&1)

if ! [[ "$config_test_result" == *"syntax check OK"* ]]; then
  bashio::log.error "Invalid configuration file"
  return 1
fi

exec dnsmasq --conf-file="$CONFIG_FILE"
