# dnsmasq-full/rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bashio
CONFIG="/etc/dnsmasq.conf"

bashio::log.info "Preparing to start..."
bashio::config.require 'data_path'

DATA_PATH=$(bashio::config 'data_path')

#Create config if not exist
if ! bashio::fs.file_exists "$DATA_PATH/dnsmasq.conf"; then
  mkdir -p "$DATA_PATH" || bashio::exit.nok "Could not create $DATA_PATH"

  cat <<EOF > "$DATA_PATH/dnsmasq.conf"
no-resolv
no-hosts
keep-in-foreground
log-facility=-
no-poll
user=root
server=8.8.8.8
server=8.8.4.4
EOF
fi

# Run dnsmasq
bashio::log.info "Starting dnsmasq..."

# Set max open file limit to speed up startup
ulimit -n 1024

exec dnsmasq -C "$DATA_PATH/dnsmasq.conf" -z < /dev/null
